# files/system/etc/systemd/system/system-flatpak-setup.service.d/override.conf
[Service]
# Disable system-flatpak-setup.timer after succesful exec
ExecStartPost=/usr/bin/bluebuild-flatpak-manager disable system
